# src/daemon/lldpd.service.in
[Unit]
Description=LLDP daemon
Documentation=man:lldpd(8)

[Service]
Type=notify
NotifyAccess=main
EnvironmentFile=-/etc/default/lldpd
ExecStart=@prefix@/sbin/lldpd $DAEMON_ARGS
ExecStartPre=mkdir -p @PRIVSEP_CHROOT@/etc
ExecStartPre=cp /etc/localtime @PRIVSEP_CHROOT@/etc/localtime
Restart=on-failure

[Install]
WantedBy=multi-user.target
