# debian/lldpd.upstart

# lldpd - LLDP daemon
#
# lldpd is a 802.1AB implementation, a L2 network discovery
# protocol. It also supports CDP, EDP and various other protocols.

description "LLDP daemon"

start on net-device-up IFACE=lo
stop on runlevel [06]

expect daemon
respawn

script
  . /etc/default/lldpd
  exec lldpd $DAEMON_ARGS
end script
